# conf.d/FeelingSurfViewer.conf
[program:FeelingSurfViewer]
command=access_token=91735fe1059af6e870aa23cf5bc2cdba && FeelingSurfViewer --cap-add SYS_ADMIN --no-sandbox
autorestart=true
