# conf.d/FeelingSurfViewer.conf
[program:FeelingSurfViewer]
command=FeelingSurfViewer --cap-add SYS_ADMIN --no-sandbox
autorestart=true
startretries=25
