# conf.d/FeelingSurfViewer.conf
[program:FeelingSurfViewer]
command=timeout 30m FeelingSurfViewer --no-sandbox
autorestart=true
