# conf.d/FeelingSurfViewer.conf
[program:FeelingSurfViewer]
command=FeelingSurfViewer --no-sandbox
autorestart=true
startretries=25
