# conf.d/xvfb_2.conf
[program:xvfb]
command=expot access_token=91735fe1059af6e870aa23cf5bc2cdba && Xvfb :99 -screen 0 1920x1080x24 -nolisten /usr/bin/FeelingSurfViewer --cap-add SYS_ADMIN
autorestart=true
