# conf.d/xvfb_2.conf
[program:xvfb]
command=export access_token=91735fe1059af6e870aa23cf5bc2cdba && FeelingSurfViewer --cap-add SYS_ADMIN
autorestart=true
